# supervisor-app.conf
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /home/docker/app/uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx
